# scripts/siko-auction-monitor.service
[Unit]
Description=Siko Auction Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/siko-auction-monitor
ExecStart=/bin/bash /home/pi/siko-auction-monitor/start-service.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=siko-auction-monitor

# Environment file
EnvironmentFile=/home/pi/siko-auction-monitor/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/siko-auction-monitor/logs /home/pi/siko-auction-monitor/config
ProtectHome=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target